# output/file_event_win_susp_task_write.json
{
  "metadata": {
    "name": "Suspicious Scheduled Task Write to System32 Tasks",
    "comment": "Detects the creation of tasks from processes executed from suspicious locations",
    "annotations": {
      "source": "sigma",
      "category": "file_event",
      "product": "windows",
      "sigma_status": "test"
    }
  },
  "spec": {
    "schedule": {
      "atLeastEvery": "12h",
      "enabled": false
    },
    "input": {
      "batch": {
        "sql": "SELECT * FROM file_activity WHERE time BETWEEN CURRENT_TIMESTAMP() - INTERVAL 24 HOUR AND CURRENT_TIMESTAMP() AND (activity_id = 3 AND contains(lower(file.path), lower('\\\\\\\\windows\\\\\\\\system32\\\\\\\\tasks')) AND (contains(lower(actor.process.name), lower('\\\\\\\\appdata\\\\\\\\')) OR contains(lower(actor.process.name), lower('c:\\\\\\\\perflogs')) OR contains(lower(actor.process.name), lower('\\\\\\\\windows\\\\\\\\system32\\\\\\\\config\\\\\\\\systemprofile'))))"
      }
    },
    "output": {
      "summary": "Detects the creation of tasks from processes executed from suspicious locations",
      "defaultContext": true
    },
    "metadata": {
      "version": 1.0,
      "severity": "High",
      "fidelity": "Investigative",
      "objective": "False positive: Unknown Reference: Internal Research",
      "mitre": [
        {
          "taxonomy": "MITRE ATT&CK",
          "tactic": "Persistence",
          "techniqueId": "T1053"
        }
      ]
    },
    "computeMode": "high"
  },
  "apiVersion": "v1",
  "kind": "Rule"
}