# output/file_event_lnx_persistence_cron_files.json
{
  "metadata": {
    "name": "Persistence Via Cron Files",
    "comment": "Detects creation of cron file or files in Cron directories which could indicates potential persistence.",
    "annotations": {
      "source": "sigma",
      "category": "file_event",
      "product": "linux",
      "sigma_status": "test"
    }
  },
  "spec": {
    "schedule": {
      "atLeastEvery": "12h",
      "enabled": false
    },
    "input": {
      "batch": {
        "sql": "SELECT * FROM file_activity WHERE time BETWEEN CURRENT_TIMESTAMP() - INTERVAL 24 HOUR AND CURRENT_TIMESTAMP() AND (activity_id = 1 AND (startswith(lower(file.path), lower('/etc/cron.d/')) OR startswith(lower(file.path), lower('/etc/cron.daily/')) OR startswith(lower(file.path), lower('/etc/cron.hourly/')) OR startswith(lower(file.path), lower('/etc/cron.monthly/')) OR startswith(lower(file.path), lower('/etc/cron.weekly/')) OR startswith(lower(file.path), lower('/var/spool/cron/crontabs/')) OR contains(lower(file.path), lower('/etc/cron.allow')) OR contains(lower(file.path), lower('/etc/cron.deny')) OR contains(lower(file.path), lower('/etc/crontab'))))"
      }
    },
    "output": {
      "summary": "Detects creation of cron file or files in Cron directories which could indicates potential persistence.",
      "defaultContext": true
    },
    "metadata": {
      "version": 1.0,
      "severity": "Medium",
      "fidelity": "Investigative",
      "objective": "False positive: Any legitimate cron file. Reference: https://github.com/microsoft/MSTIC-Sysmon/blob/f1477c0512b0747c1455283069c21faec758e29d/linux/configs/attack-based/persistence/T1053.003_Cron_Activity.xml",
      "mitre": [
        {
          "taxonomy": "MITRE ATT&CK",
          "tactic": "Persistence",
          "techniqueId": "T1053",
          "subTechniqueId": "003"
        }
      ]
    },
    "computeMode": "high"
  },
  "apiVersion": "v1",
  "kind": "Rule"
}