# output/proc_creation_win_schtasks_schedule_type_system.json
{
  "metadata": {
    "name": "Suspicious Schtasks Schedule Type With High Privileges",
    "comment": "Detects scheduled task creations or modification to be run with high privileges on a suspicious schedule type",
    "annotations": {
      "source": "sigma",
      "category": "process_creation",
      "product": "windows",
      "sigma_status": "test"
    }
  },
  "spec": {
    "schedule": {
      "atLeastEvery": "12h",
      "enabled": false
    },
    "input": {
      "batch": {
        "sql": "SELECT * FROM process_activity WHERE time BETWEEN CURRENT_TIMESTAMP() - INTERVAL 24 HOUR AND CURRENT_TIMESTAMP() AND (activity_id = 1 AND endswith(lower(process.name), lower('\\\\\\\\schtasks.exe')) AND lower(process.file.name) = lower('schtasks.exe') AND (contains(lower(process.cmd_line), lower(' onlogon ')) OR contains(lower(process.cmd_line), lower(' onstart ')) OR contains(lower(process.cmd_line), lower(' once ')) OR contains(lower(process.cmd_line), lower(' onidle '))) AND (contains(lower(process.cmd_line), lower('nt aut')) OR contains(lower(process.cmd_line), lower(' system')) OR contains(lower(process.cmd_line), lower('highest'))))"
      }
    },
    "output": {
      "summary": "Detects scheduled task creations or modification to be run with high privileges on a suspicious schedule type",
      "defaultContext": true
    },
    "metadata": {
      "version": 1.0,
      "severity": "Medium",
      "fidelity": "Investigative",
      "objective": "False positive: Some installers were seen using this method of creation unfortunately. Filter them in your environment References: https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-change, https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create",
      "mitre": [
        {
          "taxonomy": "MITRE ATT&CK",
          "tactic": "Execution",
          "techniqueId": "T1053",
          "subTechniqueId": "005"
        }
      ]
    },
    "computeMode": "high"
  },
  "apiVersion": "v1",
  "kind": "Rule"
}