# output/proc_creation_win_schtasks_persistence_windows_telemetry.json
{
  "metadata": {
    "name": "Potential Persistence Via Microsoft Compatibility Appraiser",
    "comment": "Detects manual execution of the \"Microsoft Compatibility Appraiser\" task via schtasks. In order to trigger persistence stored in the \"\\AppCompatFlags\\TelemetryController\" registry key.",
    "annotations": {
      "source": "sigma",
      "category": "process_creation",
      "product": "windows",
      "sigma_status": "test"
    }
  },
  "spec": {
    "schedule": {
      "atLeastEvery": "12h",
      "enabled": false
    },
    "input": {
      "batch": {
        "sql": "SELECT * FROM process_activity WHERE time BETWEEN CURRENT_TIMESTAMP() - INTERVAL 24 HOUR AND CURRENT_TIMESTAMP() AND (activity_id = 1 AND endswith(lower(process.name), lower('\\\\\\\\schtasks.exe')) AND lower(process.file.name) = lower('schtasks.exe') AND contains(lower(process.cmd_line), lower('run ')) AND contains(lower(process.cmd_line), lower('\\\\\\\\application experience\\\\\\\\microsoft compatibility appraiser')))"
      }
    },
    "output": {
      "summary": "Detects manual execution of the \"Microsoft Compatibility Appraiser\" task via schtasks. In order to trigger persistence stored in the \"\\AppCompatFlags\\TelemetryController\" registry key.",
      "defaultContext": true
    },
    "metadata": {
      "version": 1.0,
      "severity": "Medium",
      "fidelity": "Investigative",
      "objective": "False positive: Unknown Reference: https://www.trustedsec.com/blog/abusing-windows-telemetry-for-persistence/",
      "mitre": [
        {
          "taxonomy": "MITRE ATT&CK",
          "tactic": "Persistence",
          "techniqueId": "T1053",
          "subTechniqueId": "005"
        }
      ]
    },
    "computeMode": "high"
  },
  "apiVersion": "v1",
  "kind": "Rule"
}