# output/proc_creation_win_schtasks_system.json
{
  "metadata": {
    "name": "Schtasks Creation Or Modification With SYSTEM Privileges",
    "comment": "Detects the creation or update of a scheduled task to run with \"NT AUTHORITY\\SYSTEM\" privileges",
    "annotations": {
      "source": "sigma",
      "category": "process_creation",
      "product": "windows",
      "sigma_status": "test"
    }
  },
  "spec": {
    "schedule": {
      "atLeastEvery": "12h",
      "enabled": false
    },
    "input": {
      "batch": {
        "sql": "SELECT * FROM process_activity WHERE time BETWEEN CURRENT_TIMESTAMP() - INTERVAL 24 HOUR AND CURRENT_TIMESTAMP() AND (lower(activity_id) = lower('1') AND endswith(lower(process.name), lower('\\\\schtasks.exe')) AND (contains(lower(process.cmd_line), lower(' /change ')) OR contains(lower(process.cmd_line), lower(' /create '))) AND contains(lower(process.cmd_line), lower('/ru ')) AND (contains(lower(process.cmd_line), lower('nt aut')) OR contains(lower(process.cmd_line), lower(' system '))) AND NOT (endswith(lower(process.name), lower('\\\\schtasks.exe')) AND contains(lower(process.cmd_line), lower('/tn tvinstallrestore')) AND contains(lower(process.cmd_line), lower('\\\\teamviewer_.exe')) OR contains(lower(process.cmd_line), lower('/create /f /ru system /sc weekly /tn avirasystemspeedupverify /tr ')) OR contains(lower(process.cmd_line), lower(':\\\\program files (x86)\\\\avira\\\\system speedup\\\\setup\\\\avira_speedup_setup.exe')) OR contains(lower(process.cmd_line), lower('/verify /verysilent /nostart /nodotnet /norestart\" /rl highest'))))"
      }
    },
    "output": {
      "summary": "Detects the creation or update of a scheduled task to run with \"NT AUTHORITY\\SYSTEM\" privileges",
      "defaultContext": true
    },
    "metadata": {
      "version": 1.0,
      "severity": "High",
      "fidelity": "Investigative",
      "objective": "False positive: Unknown References: https://www.elastic.co/security-labs/exploring-the-qbot-attack-pattern, https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks",
      "mitre": [
        {
          "taxonomy": "MITRE ATT&CK",
          "tactic": "Execution",
          "techniqueId": "T1053",
          "subTechniqueId": "005"
        }
      ]
    },
    "computeMode": "high"
  },
  "apiVersion": "v1",
  "kind": "Rule"
}