# output/win_security_susp_scheduled_task_update.json
{
  "metadata": {
    "name": "Suspicious Scheduled Task Update",
    "comment": "Detects update to a scheduled task event that contain suspicious keywords.",
    "annotations": {
      "source": "sigma",
      "product": "windows",
      "service": "security",
      "sigma_status": "test"
    }
  },
  "spec": {
    "schedule": {
      "atLeastEvery": "15m",
      "enabled": true
    },
    "input": {
      "batch": {
        "sql": "SELECT * FROM scheduled_job_activity WHERE time BETWEEN CURRENT_TIMESTAMP() - INTERVAL 24 HOUR AND CURRENT_TIMESTAMP() AND (lower(activity_id) = lower('2') AND metadata.event_code = 4702 AND (contains(lower(job.cmd_line), lower('\\\\appdata\\\\local\\\\temp\\\\')) OR contains(lower(job.cmd_line), lower('\\\\appdata\\\\roaming\\\\')) OR contains(lower(job.cmd_line), lower('\\\\users\\\\public\\\\')) OR contains(lower(job.cmd_line), lower('\\\\windows\\\\temp\\\\')) OR contains(lower(job.cmd_line), lower('c:\\\\temp\\\\')) OR contains(lower(job.cmd_line), lower('\\\\desktop\\\\')) OR contains(lower(job.cmd_line), lower('\\\\downloads\\\\')) OR contains(lower(job.cmd_line), lower('\\\\temporary internet')) OR contains(lower(job.cmd_line), lower('c:\\\\programdata\\\\')) OR contains(lower(job.cmd_line), lower('c:\\\\perflogs\\\\'))) AND (contains(lower(job.cmd_line), lower('regsvr32')) OR contains(lower(job.cmd_line), lower('rundll32')) OR contains(lower(job.cmd_line), lower('cmd.exe</command>')) OR contains(lower(job.cmd_line), lower('cmd</command>')) OR contains(lower(job.cmd_line), lower('<arguments>/c ')) OR contains(lower(job.cmd_line), lower('<arguments>/k ')) OR contains(lower(job.cmd_line), lower('<arguments>/r ')) OR contains(lower(job.cmd_line), lower('powershell')) OR contains(lower(job.cmd_line), lower('pwsh')) OR contains(lower(job.cmd_line), lower('mshta')) OR contains(lower(job.cmd_line), lower('wscript')) OR contains(lower(job.cmd_line), lower('cscript')) OR contains(lower(job.cmd_line), lower('certutil')) OR contains(lower(job.cmd_line), lower('bitsadmin')) OR contains(lower(job.cmd_line), lower('bash.exe')) OR contains(lower(job.cmd_line), lower('bash ')) OR contains(lower(job.cmd_line), lower('scrcons')) OR contains(lower(job.cmd_line), lower('wmic ')) OR contains(lower(job.cmd_line), lower('wmic.exe')) OR contains(lower(job.cmd_line), lower('forfiles')) OR contains(lower(job.cmd_line), lower('scriptrunner')) OR contains(lower(job.cmd_line), lower('hh.exe'))))"
      }
    },
    "output": {
      "summary": "Detects update to a scheduled task event that contain suspicious keywords.",
      "defaultContext": true
    },
    "metadata": {
      "version": 1.0,
      "severity": "High",
      "fidelity": "Investigative",
      "objective": "False positive: Unknown Reference: https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4698",
      "mitre": [
        {
          "taxonomy": "MITRE ATT&CK",
          "tactic": "Execution",
          "techniqueId": "T1053",
          "subTechniqueId": "005"
        }
      ]
    },
    "computeMode": "high"
  },
  "apiVersion": "v1",
  "kind": "Rule"
}